# infra/vps/systemd/myboon-publisher.service
[Unit]
Description=MYBOON Publisher Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/myboon/packages/brain
ExecStart=/usr/bin/env pnpm publisher
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
