# infra/vps/systemd/myboon-v3-wallet-profiles.service
[Unit]
Description=MYBOON V3 Hyperliquid Wallet Profile Snapshot
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/myboon/packages/collectors
ExecStart=/usr/bin/flock -n -E 0 /tmp/myboon-v3-wallet-profiles.lock /usr/bin/env pnpm hyperliquid:wallet-profiles
Environment=NODE_ENV=production

