# infra/vps/systemd/myboon-analyst.service
[Unit]
Description=MYBOON Narrative Analyst Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/myboon/packages/brain
ExecStart=/usr/bin/env pnpm narrative:analyst
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
